# debian/spreedbox-fpgaprog.upstart
# fpgaprog - Program FPGA
#   Programs the FPGA via SPI.

description "initialize Spreedbox FPGA"

start on (startup and started kmod)

task

console log

env DEVICE=/dev/spidev0.0
env BITMAP=/usr/share/spreedbox-fpgaprog/bitmaps/default.bin

script
	/usr/sbin/spreedbox-fpgaprog -d $DEVICE $BITMAP || true

	initctl emit --no-wait spreedbox-fpga-initialized
end script
